# config/google-calendar-mcp.service
[Unit]
Description=Google Calendar MCP Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=richard
Group=richard
WorkingDirectory=***REMOVED***/git/google_calendar_mcp
ExecStart=***REMOVED***/git/google_calendar_mcp/scripts/start_mcp_remote.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=google-calendar-mcp

Environment=PYTHONPATH=***REMOVED***/git/google_calendar_mcp/src
Environment=MCP_SERVER_HOST=0.0.0.0
Environment=MCP_SERVER_PORT=3001
Environment=MCP_TRANSPORT=sse

[Install]
WantedBy=multi-user.target 