# config/google-calendar-mcp.service
[Unit]
Description=Google Calendar MCP Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=richard
Group=richard
WorkingDirectory=***REMOVED***/git/google_calendar_mcp
ExecStart=***REMOVED***/.local/bin/uvx --from ***REMOVED***/git/google_calendar_mcp google-calendar-mcp
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=google-calendar-mcp

[Install]
WantedBy=multi-user.target 